# deploy/homelab-backend-dev.service
[Unit]
Description=HomeServerManager backend (dev VM)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
User=homelab
Group=homelab
SupplementaryGroups=docker
WorkingDirectory=/srv/homeservermanager-dev/backend
EnvironmentFile=/srv/homeservermanager-dev/backend/.env
ExecStart=/usr/bin/node /srv/homeservermanager-dev/backend/dist/server.js
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/homeservermanager

[Install]
WantedBy=multi-user.target
