# deploy/homelab-frontend-dev.service
[Unit]
Description=HomeServerManager frontend preview (dev VM)
After=network-online.target homeservermanager-backend-dev.service
Wants=network-online.target

[Service]
Type=simple
User=homelab
Group=homelab
WorkingDirectory=/srv/homeservermanager-dev/frontend
ExecStart=/usr/bin/node /srv/homeservermanager-dev/frontend/node_modules/vite/bin/vite.js preview --host 0.0.0.0 --port 4173
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict

[Install]
WantedBy=multi-user.target
